# example_demon.service
[Unit]
Description=Chatterbox TTS API Server
After=network.target

[Service]
User=ladvien
Group=ladvien
WorkingDirectory=/home/ladvien/speech_server
ExecStart=/home/ladvien/speech_server/.venv/bin/python server.py
Restart=on-failure
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
